# systemd/raspi-hub.service.example
[Unit]
Description=raspi-esp32-status-panel hub server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=<USER>
WorkingDirectory=/home/<USER>/raspi-esp32-status-panel/hub
ExecStart=/home/<USER>/raspi-esp32-status-panel/hub/venv/bin/python3 /home/<USER>/raspi-esp32-status-panel/hub/main.py
Restart=always
RestartSec=3
Environment=HUB_DATA_DIR=/home/<USER>/raspi-esp32-status-panel/hub/data
Environment=HUB_PORT=8080

[Install]
WantedBy=multi-user.target
